# kfrgb.service
[Unit]
Description=Set Kingston Fury RAM RGB 

[Service]
ExecStart=/opt/kfrgb/i2cset-static-blue.sh

[Install]
WantedBy=multi-user.target
